# tp3/prueba_docker/supervisord.conf
[supervisord]
nodaemon=true

[program:levantar_apache]
command=/script/apache_inicio.sh

